# systemd/user/battery_notify.timer
[Unit]
Description=Run Battery Notification Service every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
